# systemd/deckery-hud.service
[Unit]
Description=Deckery HUD overlay (persistent D-Bus service)
PartOf=graphical-session.target
PartOf=deckery-tray.service
After=graphical-session.target
After=deckery-tray.service

[Service]
Type=dbus
BusName=de.plasma_deckery.hud
ExecStart=%h/.local/bin/deckery-hud
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
